# deploy/all-in-one/supervisord.conf
[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[program:mysql]
command=/usr/sbin/mysqld --user=mysql --datadir=%(ENV_ASR_MYSQL_DATA_DIR)s --socket=/run/mysqld/mysqld.sock --pid-file=/run/mysqld/mysqld.pid --bind-address=127.0.0.1 --port=3306
priority=10
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:admin-api]
command=/bin/sh -lc 'export TMPDIR="$ASR_TMP_DIR"; /app/scripts/wait-for-mysql.sh && cd /app/backend && /app/bin/admin-api'
priority=20
autostart=true
autorestart=true
startretries=20
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:asr-api]
command=/bin/sh -lc 'export TMPDIR="$ASR_TMP_DIR"; /app/scripts/wait-for-mysql.sh && cd /app/backend && /app/bin/asr-api'
priority=20
autostart=true
autorestart=true
startretries=20
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nlp-api]
command=/bin/sh -lc 'export TMPDIR="$ASR_TMP_DIR"; /app/scripts/wait-for-mysql.sh && cd /app/backend && /app/bin/nlp-api'
priority=20
autostart=true
autorestart=true
startretries=20
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gateway]
command=/bin/sh -lc 'cd /app/backend && /app/bin/gateway'
priority=30
autostart=true
autorestart=true
startretries=10
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
priority=40
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0